# systemd/pricebuzz.service
[Unit]
Description=PriceBuzz
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/YOUR_USER/pricebuzz
ExecStart=/home/YOUR_USER/.nvm/versions/node/v22.22.2/bin/node dist/server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
EnvironmentFile=/home/YOUR_USER/pricebuzz/.env

[Install]
WantedBy=multi-user.target
